# deploy/systemd/kirokiro.service
[Unit]
Description=Kirokiro Gunicorn Application Server
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/home/kirokiro/kirokiro
Environment="PATH=/home/kirokiro/kirokiro/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=kirokiro.settings"
EnvironmentFile=/home/kirokiro/kirokiro/.env
ExecStart=/home/kirokiro/kirokiro/venv/bin/gunicorn -c /home/kirokiro/kirokiro/deploy/gunicorn/gunicorn.conf.py kirokiro.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/kirokiro/kirokiro/media /home/kirokiro/kirokiro/staticfiles /run/gunicorn /var/log/gunicorn /home/kirokiro/kirokiro/logs

[Install]
WantedBy=multi-user.target
